# deploy/minilpa-backend.service
[Unit]
Description=MiniLPA Web Backend Service
After=network.target

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/minilpa/app
ExecStart=/usr/bin/java -Xms256m -Xmx512m -jar /www/wwwroot/minilpa/app/minilpa-backend.jar --spring.config.location=/www/wwwroot/minilpa/config/application.yml
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量
Environment="JAVA_HOME=/usr/lib/jvm/java-21-openjdk"  # 根据实际Java路径修改
Environment="SPRING_PROFILES_ACTIVE=prod"

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

